# podman/quadlets/router-ai-atius-redis.container
# systemd quadlet for Redis (router-ai-atius-redis)
#
# Rebrand v2.11: was redis-newapi, now router-ai-atius-redis.

[Unit]
Description=router-ai-atius-redis (Redis 7 — cache + rate limiting)
After=network-online.target

[Container]
Image=docker.io/library/redis:7-alpine
ContainerName=router-ai-atius-redis
# Password is injected via EnvironmentFile in production. Dev default below.
Environment=REDIS_PASSWORD=change-me-in-production
Command=redis-server --requirepass ${REDIS_PASSWORD}
Timezone=local
NoNewPrivileges=true
ReadOnly=true
Tmpfs=/data
HealthCmd=redis-cli -a ${REDIS_PASSWORD} ping
HealthIntervalSec=30s
HealthTimeoutSec=5s
HealthRetries=3

[Service]
Restart=always
RestartSec=10
User=%U
Group=%U

[Install]
WantedBy=default.target
